# systemd/audiocinema.service
[Unit]␊
Description=AudioCinema scheduled automation run
After=default.target␊
␊
[Service]␊
Type=oneshot␊
WorkingDirectory=%h/code␊
ExecStart=%h/code/venv/bin/python %h/code/src/main.py --autorun
